# etc/services.d/cron/run
#!/usr/bin/with-contenv bash

if [ -f /usr/sbin/cron ]; then # ubuntu
  exec /usr/sbin/cron -f
else # alpine
  exec /usr/sbin/crond -n
fi
