# etc/services.d/rclone/finish
#!/usr/bin/with-contenv bash

/bin/fusermount -uz /gdrive-cloud/*