# etc/services.d/mergerfs/finish
#!/usr/bin/with-contenv bash

touch /config/secretfile

/bin/fusermount -uz /remote